# config/intercomd.service
[Unit]
Description=Intercom Rust daemon (IPC bridge, event consumer, Telegram bridge)
After=network-online.target
Wants=network-online.target
# Start before intercom.service so IPC queries are handled from boot
Before=intercom.service

[Service]
Type=simple
WorkingDirectory=/home/mk/projects/Demarch/apps/intercom
ExecStart=/home/mk/projects/Demarch/apps/intercom/rust/target/release/intercomd serve --config config/intercom.toml
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Graceful shutdown — give IPC watcher time to finish current poll cycle
TimeoutStopSec=10

# Environment overrides (uncomment to customize)
# Environment=INTERCOMD_BIND=127.0.0.1:7340
# Environment=HOST_CALLBACK_URL=http://127.0.0.1:7341
# Environment=RUST_LOG=intercomd=info

[Install]
WantedBy=default.target
